# roles/borgbackup-server/templates/update-motd.d/05-service.j2
#!/usr/bin/tail +14
{{ ansible_header | comment }}
[0m> [38;5;82mBorgbackup (Serveur)[0m a été déployé sur cette machine. Les backups sont situés dans [38;5;6m{{ borg.path }}[0m.
